# services/htop
#! /bin/sh

### BEGIN INIT INFO
# Provides:          ttyd-htop
# Required-Start:    $network $remote_fs
# Required-Stop:     $network $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: ttyd htop service
### END INIT INFO

set -e

PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

NAME="ttyd-htop"
DAEMON="/usr/bin/ttyd"
PORT="8027"
PIDFILE="/run/ttyd_htop.pid"
LOGFILE="/var/log/ttyd_htop.log"

. /lib/lsb/init-functions

start() {
    log_daemon_msg "Starting ttyd htop" "$NAME"

    start-stop-daemon --start --background --make-pidfile \
        --pidfile "$PIDFILE" \
        --exec "$DAEMON" -- \
        -p "$PORT" htop >> "$LOGFILE" 2>&1

    log_end_msg $?
}

stop() {
    log_daemon_msg "Stopping ttyd htop" "$NAME"

    start-stop-daemon --stop --quiet --pidfile "$PIDFILE" --retry 30 || true
    rm -f "$PIDFILE"

    log_end_msg $?
}

restart() {
    stop
    start
}

status() {
    if [ -f "$PIDFILE" ] && kill -0 "$(cat $PIDFILE)" 2>/dev/null; then
        echo "ttyd htop is running on port $PORT"
        exit 0
    else
        echo "ttyd htop is not running"
        exit 1
    fi
}

case "$1" in
    start) start ;;
    stop) stop ;;
    restart) restart ;;
    status) status ;;
    *) echo "Usage: $0 {start|stop|restart|status}" ; exit 1 ;;
esac

exit 0