# services/ttyd
#!/bin/bash

### BEGIN INIT INFO
# Provides:          ttyd-tmux-terminal
# Required-Start:    $network $remote_fs
# Required-Stop:     $network $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: ttyd tmux terminal session
### END INIT INFO

PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

NAME="ttyd-tmux-terminal"
DAEMON="/usr/bin/ttyd"
PORT="8026"
SESSION="terminal"
WORKDIR="/home/ubuntu/home"
TMUX_SOCKET="/tmp/tmux-main"
PIDFILE="/run/ttyd_tmux_terminal.pid"
LOGFILE="/var/log/ttyd_tmux_terminal.log"
USER="ubuntu"

. /lib/lsb/init-functions

start() {
    log_daemon_msg "Starting ttyd tmux" "$NAME"

    touch "$LOGFILE"
    chown "$USER" "$LOGFILE"

    if [[ -f "$PIDFILE" ]] && kill -0 "$(cat "$PIDFILE")" 2>/dev/null; then
        echo "$NAME is already running"
        log_end_msg 0
        return 0
    fi

    sudo su -s /bin/bash "$USER" -c "
        tmux -S '$TMUX_SOCKET' new-session -d -s '$SESSION' -c '$WORKDIR' 2>/dev/null || true
        chmod 777 '$TMUX_SOCKET'
        nohup $DAEMON -p $PORT -W tmux -S '$TMUX_SOCKET' attach-session -t '$SESSION' >> '$LOGFILE' 2>&1 &
        echo \$! > '$PIDFILE'
    "

    sleep 1

    if [[ -f "$PIDFILE" ]] && kill -0 "$(cat "$PIDFILE")" 2>/dev/null; then
        log_end_msg 0
    else
        log_end_msg 1
    fi
}

stop() {
    log_daemon_msg "Stopping ttyd tmux" "$NAME"

    if [[ -f "$PIDFILE" ]]; then
        kill "$(cat "$PIDFILE")" 2>/dev/null || true
        rm -f "$PIDFILE"
    fi

    sudo su -s /bin/bash "$USER" -c "tmux -S '$TMUX_SOCKET' kill-session -t '$SESSION' 2>/dev/null || true"

    log_end_msg 0
}

restart() {
    stop
    sleep 1
    start
}

status() {
    if [[ -f "$PIDFILE" ]] && kill -0 "$(cat "$PIDFILE")" 2>/dev/null; then
        echo "$NAME is running (pid $(cat "$PIDFILE")) on port $PORT"
        exit 0
    else
        echo "$NAME is not running"
        exit 1
    fi
}

case "$1" in
    start)   start   ;;
    stop)    stop    ;;
    restart) restart ;;
    status)  status  ;;
    *) echo "Usage: $0 {start|stop|restart|status}"; exit 1 ;;
esac

exit 0